# crontab
#!bin/bash

# change wallpaper
*/10 * * * * DISPLAY=:0 feh --randomize --bg-scale ~/Dropbox/Wallpapers/new
