# systemd/user/awww-cycle.timer
[Unit]
Description=Timer for awww-cycle-wallpaper

[Timer]
OnBootSec=1
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
