# systemd/user/swww-cycle.service
[Unit]
Description=Cycle wallpapers randomly via swww

[Service]
Type=oneshot
ExecStart=/home/zander/dotfiles/scripts/swww-cycle-wallpaper.sh
